# docker/supervisor.d/zz-laravel-boot-order.conf
; php-fpm nginxdan oldin ishga tushishi kerak — aks holda birinchi so'rovlar 502 (fastcgi 9000 refused).
; Supervisord: kichik priority = yuqori ustuvorlik (guruh avval start qilinadi).
; Asl webdevops fayllarida ikkala guruh ham priority=20 bo'lib, nginx ko'pincha birinchi bo'lib ochiladi.

[group:php-fpm]
priority=5

[group:nginx]
priority=40

; Supervisord 4.2+: nginx php-fpmd tayyor bo'lgach ishga tushsin (ustuvorlik bilan birga).
[program:nginxd]
depends_on=php-fpmd
